# etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
auto_start=true
autorestart=true